# apache2_http/service/apache2/run
#!/bin/sh
exec apache2ctl start>>/var/log/apache2.log 2>&1
